# polling-311.conf
description "polling 311"
author "Eunice Kokor"
start on filesystem and net-device-up
stop on runlevel [06]

exec /home/ubuntu/polling-db/311.py

respawn
